# loopholes/claude-oauth-broker/claude-oauth-broker.service
[Unit]
Description=Claude OAuth broker for yolo-jail (MITM refresh serializer)
Documentation=https://github.com/mschulkind-oss/yolo-jail/blob/main/docs/claude-oauth-mitm-proxy-plan.md
After=network-online.target
Wants=network-online.target

[Service]
# Placeholder replaced by `just deploy` with the absolute path of the
# broker console script shipped by the yolo-jail wheel.
#
# Bind 0.0.0.0 — 169.254.1.2 is a container-side alias (pasta tunnel),
# not a real host interface, so bind(2) against it fails with
# EADDRNOTAVAIL and the daemon crash-loops. The kernel routes
# container-origin packets for `platform.claude.com` to the host's
# loopback via the --add-host :host-gateway line the yolo loopholes
# loader emits, so 0.0.0.0:443 catches them.
ExecStart=@@BROKER_BIN@@ --host 0.0.0.0 --port 443
Restart=on-failure
RestartSec=5s

# Privilege: listening on :443 needs CAP_NET_BIND_SERVICE. User-level systemd
# can grant this via AmbientCapabilities when NoNewPrivileges is off and the
# kernel allows it. If this fails on your host, switch --port to 8443 and
# arrange DNAT on the container bridge (see loopholes/claude-oauth-broker/README).
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Sandboxing — the daemon only needs to read/write the loophole dir and the
# shared credentials file.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/yolo-jail/loopholes/claude-oauth-broker %h/.local/share/yolo-jail/home/.claude
PrivateTmp=yes
NoNewPrivileges=no

[Install]
WantedBy=default.target
