# scripts/claude-token-refresher.service
[Unit]
Description=Refresh shared Claude Code OAuth token for yolo-jail
Documentation=https://github.com/mschulkind-oss/yolo-jail/blob/main/scripts/README.md
# Don't run before network is up — we need to reach platform.claude.com
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# ExecStart's placeholder token below is sed-replaced by `just deploy` with
# the absolute path to the claude-token-refresher binary that uv tool
# install / pip install shipped as part of the yolo-jail wheel.  No repo
# checkout required on the host at runtime.
ExecStart=@@REFRESHER_BIN@@
# Don't retry on failure — the timer will fire again soon enough.  A
# transient network error shouldn't tie up systemd.
Restart=no
# Defensive sandboxing.  Most of /home is read-only; only the directories
# the refresher actually writes are carved out:
#   - ~/.claude                  (the host credentials file)
#   - ~/.local/share/yolo-jail   (mirrored creds + the refresh.lock flock)
# ProtectHome=read-write is NOT a valid systemd value — systemd silently
# ignored the directive and left HOME fully writable.  This replaces it
# with the real thing.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.claude %h/.local/share/yolo-jail
PrivateTmp=yes
NoNewPrivileges=yes
# Timeout: generous, in case the auth endpoint is slow.
TimeoutStartSec=60
